# service-ag-master/bin/Debug/netcoreapp2.2/instasoft.service
[Unit]
Description=Instasoft ASP.NET WebApi server
After=multi-user.target
Requires=mysql.service

[Service]
WorkingDirectory=/var/www/instasoft/publish
ExecStart=sudo /var/www/instasoft/publish/instasoft
Restart=always

[Install]
WantedBy=multi-user.target

